# test1/crontab
* * * * * root python3 /usr/local/bin/monitor.py >> /var/log/cron.log 2>&1
